# systemd/ai-assistant-api.service
[Unit]
Description=AI Assistant REST API
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ai-assistant
Group=ai-assistant
WorkingDirectory=/opt/ai-assistant
Environment="PATH=/opt/ai-assistant/venv/bin"
ExecStart=/opt/ai-assistant/venv/bin/uvicorn api.app:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=10
StandardOutput=append:/var/log/ai-assistant/api.log
StandardError=append:/var/log/ai-assistant/api-error.log

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ai-assistant/logs /opt/ai-assistant/data /opt/ai-assistant/cache
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65535
MemoryLimit=2G
CPUQuota=100%

[Install]
WantedBy=multi-user.target
